# Install/Lwm2mClient-0.0.1/lwm2mclient.service
[Unit]
Description=Some text

[Service]
ExecStart=/opt/LWM2MClient/Start.sh

[Install]
WantedBy=multi-user.target
